# powersaving/powersaving-resume.service
[Unit]
Description=Local system resume actions
After=suspend.target

[Service]
Type=simple
ExecStart=/usr/bin/systemctl restart powersaving.service

[Install]
WantedBy=suspend.target